# conf/centos7/cron.daily/certbot
#!/bin/sh

/usr/bin/certbot renew -n --post-hook="systemctl reload httpd"
